# dist/obsworker.service
[Unit]
Description=OBS worker daemon
After=syslog.target network.target

[Service]
ExecStart=/usr/libexec/obsworker start
ExecStop=/usr/libexec/obsworker stop
Type=forking

[Install]
WantedBy=multi-user.target
