# dist/fedora/obsdispatcher.service
[Unit]
Description=OBS dispatcher daemon
After=syslog.target

[Service]
ExecStart=/usr/libexec/obsdispatcher start
Type=forking

[Install]
WantedBy=multi-user.target
