# dist/fedora/obsrepserver.service
[Unit]
Description=OBS worker daemon
After=syslog.target

[Service]
ExecStart=/usr/libexec/obsrepserver start
ExecStop=/usr/libexec/obsrepserver stop
Type=forking

[Install]
WantedBy=multi-user.target
