# dist/fedora/obsservice.service
[Unit]
Description=OBS service daemon
After=syslog.target

[Service]
ExecStart=/usr/libexec/obsservice start
ExecStop=/usr/libexec/obsservice stop
Type=forking

[Install]
WantedBy=multi-user.target
