# dist/fedora/obssrcserver.service
[Unit]
Description=OBS obssrcserver daemon
After=syslog.target

[Service]
ExecStart=/usr/sbin/obssrcserver start
Type=forking

[Install]
WantedBy=multi-user.target
